# html2bash/etc/udev/rules.d/99-html2bash.rules
KERNEL=="usb_plug", ACTION=="add", RUN+="/usr/local/html2bash/check.sh"
