# autopatch-khd/etc/udev/rules.d/autopatch-khd.rules
KERNEL=="fsl-usb2-udc", ACTION=="remove", RUN+="/usr/local/AutoPatchKHD/autopatch-khd.sh"
